# guest/truenas-lab-firstboot.service
[Unit]
Description=TrueNAS lab first-boot OVF customization
After=multi-user.target vmtoolsd.service
Wants=vmtoolsd.service

[Service]
Type=oneshot
ExecStart=/opt/truenas-lab/truenas-lab-firstboot.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
